# Instal/services.d/api/run
#!/usr/bin/execlineb -P
su ${BSM_USER} -c 'R --no-save -f Tableros/Explorador_IVs_API_run.R'
