# Instal/services.d/healthcheck/run
#!/usr/bin/execlineb -P
curl -sf http://localhost:3838/__docs__/ > /dev/null
